# srcpkgs/qubes-core-agent-linux/files/cron.hourly/qubes-sync-time
#!/bin/sh
/var/run/qubes-service/clocksync
