# srcpkgs/qubes-core-agent-linux/files/cron.hourly/qubes-update-check
#!/bin/sh

/var/run/qubes-service/qubes-update-check
