# srcpkgs/qubes-core-qubesdb/files/rc.local
#!/bin/sh
unset NOTIFY_SOCKET
mkdir -p /var/run/qubes
touch /var/log/qubes/qubesdb.dom0.log
/usr/bin/qubesdb-daemon 0
