# srcpkgs/qubes-core-agent-linux/files/rc.local
#!/bin/sh

#oneshot qubes-misc-post
/usr/lib/qubes/init/misc-post.sh
